# systemd/archmap.service
[Unit]
Description=ArchMap GeoJSON/KML generator

[Service]
Type=simple
ExecStart=/usr/bin/archmap.py
