# systemd/archmap.service
[Unit]
Description=ArchMap geojson generator

[Service]
Type=simple
ExecStart=/usr/bin/archmap.py
